# overlays/deploy/authelia.container
[Unit]
Description=Authelia SSO authentication server

[Container]
Image=ghcr.io/authelia/authelia:latest
ContainerName=authelia
Network=exousia.network:alias=authelia
PublishPort=127.0.0.1:9091:9091
Volume=%h/.config/authelia:/config:z
EnvironmentFile=%h/.config/authelia/smtp.env
Label=io.containers.autoupdate=registry

HealthCmd=wget -qO /dev/null http://localhost:9091/api/health || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=15s

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
