# overlays/deploy/coredns.container
[Unit]
Description=CoreDNS local service resolver

[Container]
Image=docker.io/coredns/coredns:1.12.1
ContainerName=coredns
AutoUpdate=registry
Network=exousia.network
PublishPort=127.0.0.1:5354:53/tcp
PublishPort=127.0.0.1:5354:53/udp
Volume=%h/.config/coredns:/etc/coredns:ro,z
Exec=-conf /etc/coredns/Corefile

HealthCmd=["true"]
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=15s

[Service]
Restart=always

[Install]
WantedBy=default.target
