# overlays/deploy/exousia-registry.container
[Unit]
Description=Local container registry for development
After=exousia-network.service
Requires=exousia-network.service

[Container]
Image=docker.io/library/registry:2
ContainerName=exousia-registry
PublishPort=127.0.0.1:5000:5000
Volume=exousia-registry-data.volume:/var/lib/registry
Network=exousia.network:alias=registry
Label=io.containers.autoupdate=registry
Environment=REGISTRY_STORAGE_DELETE_ENABLED=true

HealthCmd=["wget", "-qO-", "--spider", "http://localhost:5000/v2/"]
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=15s

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
