# overlays/deploy/forgejo-runner.container
[Unit]
Description=Forgejo Actions runner
After=forgejo.service
Requires=forgejo.service

[Container]
Image=code.forgejo.org/forgejo/runner:9.1.1
ContainerName=forgejo-runner
Volume=forgejo-runner-data.volume:/data:U
Volume=%t/podman/podman.sock:/var/run/docker.sock:ro
Network=exousia.network:alias=forgejo-runner
Label=io.containers.autoupdate=registry
UserNS=keep-id
SecurityLabelDisable=true

# Registration is a manual one-time step — see docs/forgejo-runner.md
# After registering, the runner state persists in /data/.runner
WorkingDir=/data
Exec=forgejo-runner -c /data/config.yaml daemon

HealthCmd=["cat", "/proc/1/status"]
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
