# overlays/deploy/forgejo-db.container
[Unit]
Description=Forgejo PostgreSQL database
After=exousia-network.service
Requires=exousia-network.service

[Container]
Image=docker.io/library/postgres:14-alpine
ContainerName=forgejo-db
Environment=POSTGRES_USER=forgejo
Environment=POSTGRES_PASSWORD=forgejo
Environment=POSTGRES_DB=forgejo
Volume=forgejo-db-data.volume:/var/lib/postgresql/data
Network=exousia.network:alias=forgejo-db
Label=io.containers.autoupdate=registry

HealthCmd=["pg_isready", "-U", "forgejo"]
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=15s

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
