# overlays/deploy/ollama.container
[Unit]
Description=Ollama local LLM inference server

[Container]
Image=docker.io/ollama/ollama:latest
ContainerName=ollama
PublishPort=127.0.0.1:11434:11434
Volume=ollama-data.volume:/root/.ollama
Network=exousia.network:alias=ollama
Label=io.containers.autoupdate=registry

[Service]
Restart=always
TimeoutStartSec=900
ExecStartPost=/usr/bin/podman exec ollama ollama pull qwen3:8b
ExecStartPost=/usr/bin/podman exec ollama ollama pull llama3.2:1b
ExecStartPost=/usr/bin/podman exec ollama ollama pull qwen2.5-coder:1.5b
ExecStartPost=/usr/bin/podman exec ollama ollama pull moondream:1.8b

[Install]
WantedBy=default.target
