# overlays/deploy/paperless-db.container
[Unit]
Description=Paperless-ngx PostgreSQL database

[Container]
Image=docker.io/library/postgres:18
ContainerName=paperless-db
Environment=POSTGRES_DB=paperless
Environment=POSTGRES_USER=paperless
Environment=POSTGRES_PASSWORD=paperless
Volume=paperless-db-data.volume:/var/lib/postgresql
Network=exousia.network:alias=paperless-db
Label=io.containers.autoupdate=registry

HealthCmd=pg_isready -U paperless || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=15s

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
