# overlays/deploy/temporal-server.container
[Unit]
Description=Temporal workflow engine
After=temporal-db.service
Requires=temporal-db.service

[Container]
Image=docker.io/temporalio/auto-setup:latest
ContainerName=temporal-server
PublishPort=127.0.0.1:7233:7233
Network=exousia.network:alias=temporal
Environment=DB=postgres12
Environment=DB_PORT=5432
Environment=POSTGRES_USER=temporal
Environment=POSTGRES_PWD=temporal
Environment=POSTGRES_SEEDS=temporal-db
Environment=SKIP_DYNAMIC_CONFIG_SETUP=true
Label=io.containers.autoupdate=registry

HealthCmd=wget -q -O /dev/null http://127.0.0.1:7233/ || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
