# overlays/base/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
OnBootSec=
OnUnitInactiveSec=
RandomizedDelaySec=
OnCalendar=*-*-* 23:15:00 America/New_York
Persistent=true
