# overlays/base/systemd/user/mirror-base-image.timer
[Unit]
Description=Mirror base image daily

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
